# stow_root/timers/etc/systemd/system/check-updates.service
[Unit]
Description=Check for updates

[Service]
Type=oneshot
ExecStart=/home/isac/check_updates.sh
